# roles/relay-rest-api/templates/relay-rest-api.service.j2
[Service]
ExecStart=/usr/bin/node {{install_dir}}/server.js
Restart=always
User=root

[Install]
WantedBy=multi-user.target
